# keepalived_vip/templates/etc/network/if-pre-up.d/ifcfg-lvs
#!/bin/bash

/sbin/ifconfig <%= name %> up
/sbin/ifconfig <%= name %> <%= lvs_ip  %> netmask 255.255.255.255
